# airflow-apps/supervisord.conf
[supervisord]
nodaemon=true

[program:airflow_scheduler]
command=airflow scheduler

[program:airflow_webserver]
command=airflow webserver