# airflow-apps/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:airflow_scheduler]
command=airflow scheduler

[program:airflow_webserver]
command=airflow webserver